# ESIF/Packages/Installers/chrome/dtt.conf
# Copyright (c) 2014 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Intel(R) Dynamic Tuning daemon"
author          "dptf@lists.01.org"

start on starting system-services
stop on stopping system-services

respawn

script
	exec minijail0 -u daemon -g daemon /usr/bin/ipfhostd --load:Dptf.so --srv:ipfsrv.dptf --foreground
end script
