# Packages/Installers/chrome/dptf.conf
# Copyright (c) 2014 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Intel(R) Dynamic Platform and Thermal Framework daemon"
author          "dptf@lists.01.org"

start on starting system-services
stop on stopping system-services

expect fork
respawn

exec esif_ufd
