# Packages/Installers/chrome/dptf.conf
# Copyright (c) 2014 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Start the Intel(R) Dynamic Platform and Thermal Framework daemon"
author          "dptf@lists.01.org"

start on starting system-services
stop on stopping system-services

console none

expect fork

pre-start script
  insmod /lib/modules/`uname -r`/extra/esif_lf.ko
  insmod /lib/modules/`uname -r`/extra/dptf_acpi.ko
end script

script
  exec /usr/bin/esif_ufd
end script

post-stop script
  rmmod dptf_acpi.ko
  rmmod esif_lf.ko
end script
