# ESIF/Packages/Installers/linux/dptf.service
[Unit]
description=Intel(R) Dynamic Tuning daemon

[Service]
Type=forking
Restart=always
ExecStart=/usr/bin/esif_ufd

[Install]
WantedBy=sysinit.target
