# ESIF/Packages/Installers/linux/ipf.service
[Unit]
Description=Intel(R) Innovation Platform Framework

[Service]
NoNewPrivileges=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
KeyringMode=private
RestrictRealtime=true
PrivateTmp=true
ProtectHostname=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
ProtectSystem=true
ProtectHome=true
RestrictNamespaces=true
LockPersonality=true
SystemCallArchitectures=native
CapabilityBoundingSet=CAP_SYS_RAWIO CAP_CHOWN CAP_NET_ADMIN
UMask=0077
Type=forking
Restart=always
ExecStart=/usr/bin/ipf_ufd

[Install]
WantedBy=sysinit.target
